# utils/pase-index.timer
[Unit]
Description=PaSe indexing
Requires=pase-index.service

[Timer]
Unit=pase-index.service
OnCalendar=daily

[Install]
WantedBy=timers.target
